# dot_config/systemd/user/zoekt-chromium-overlay.service
[Unit]
Description=Zoekt overlay watcher (working tree)

[Service]
Type=simple
Environment=CHROMIUM_SRC=%h/src/chromium/src
ExecStart=%h/.local/bin/zoekt-overlay-watch
Restart=always
RestartSec=2
