# dot_config/systemd/user/zoekt-chromium-index.service
[Unit]
Description=Zoekt incremental HEAD index for Chromium

[Service]
Type=oneshot
Environment=CHROMIUM_SRC=%h/src/chromium/src
ExecStart=%h/.local/bin/zoekt-index-chromium
